# root/etc/services.d/tautulli/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc python /opt/tautulli/Tautulli.py \
	--datadir /config
